# nginx/nginx.super.conf
[program:nginx]
command=/usr/sbin/nginx
# stdout_logfile=/var/log/nginx/access.log
 stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log